# config/supervision/zen-brain1-daily-sweep.timer
[Unit]
Description=zen-brain1 Daily Full Sweep Timer (6 AM daily)

[Timer]
OnCalendar=0 6 * * *
Persistent=true
AccuracySec=300

[Install]
WantedBy=timers.target
