# config/supervision/zen-brain1-l1.service
[Unit]
Description=zen-brain1 L1 Worker (llama.cpp Qwen3.5-0.8B, 10 parallel slots)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=neves
WorkingDirectory=/home/neves
ExecStart=/home/neves/git/llama.cpp/build/bin/llama-server \
    -m /home/neves/git/ai/Qwen3.5-0.8B-Q4_K_M.gguf \
    --port 56227 --host 0.0.0.0 \
    --parallel 10 --ctx-size 65536 \
    --log-disable
ExecStartPost=/bin/sleep 12
Restart=on-failure
RestartSec=10
RuntimeDirectory=zen-brain1
LogsDirectory=zen-brain1
StandardOutput=append:/var/log/zen-brain1/l1-worker.log
StandardError=append:/var/log/zen-brain1/l1-worker.log

# Resource limits
MemoryMax=8G
CPUWeight=512

[Install]
WantedBy=multi-user.target
