# config/supervision/zen-brain1-quad-hourly-summary.timer
[Unit]
Description=zen-brain1 Quad-Hourly Summary Timer (every 4 hours)

[Timer]
OnCalendar=0 */4 * * *
Persistent=true
AccuracySec=120

[Install]
WantedBy=timers.target
