# config/supervision/zen-brain1-daily-sweep.service
[Unit]
Description=zen-brain1 Daily Full Sweep (all 10 task classes)

[Service]
Type=oneshot
User=neves
Environment=BATCH_NAME=daily-full-sweep
Environment=OUTPUT_ROOT=/var/lib/zen-brain1/runs
Environment=TASKS=dead_code,defects,tech_debt,roadmap,bug_hunting,stub_hunting,package_hotspots,test_gaps,config_drift,executive_summary
Environment=TIMEOUT=300
Environment=WORKERS=5
Environment=L1_ENDPOINT=http://localhost:56227/v1/chat/completions
Environment=L1_MODEL=Qwen3.5-0.8B-Q4_K_M.gguf
ExecStart=/home/neves/zen/zen-brain1/cmd/useful-batch/useful-batch
WorkingDirectory=/home/neves/zen/zen-brain1
TimeoutStartSec=3600
StandardOutput=append:/var/log/zen-brain1/schedules.log
StandardError=append:/var/log/zen-brain1/schedules.log
